# modules/Installers/LinuxJavaFXApplicationArmV6EmbeddedJREDEB/src/deb/control/postinst
#!/bin/sh

echo "Starting [[tft.deb.name]] ... "
systemctl enable [[tft.deb.name]]
systemctl start [[tft.deb.name]]

echo "The service will be automatically started at boot."
echo "To disable it, run:"
echo " "
echo "    systemctl disable [[tft.deb.name]]"
echo " "
echo "To re-enable it, run:"
echo " "
echo "    systemctl enable [[tft.deb.name]]"
echo " "
