# modules/Installers/LinuxJavaFXApplicationArmV6EmbeddedJREDEB/src/deb/control/prerm
#!/bin/sh

echo "Stopping bluemarine2... "
systemctl stop bluemarine2
systemctl disable bluemarine2
