# modules/Installers/LinuxJavaFXApplicationArm64EmbeddedJREDEB/src/deb/control/prerm
#!/bin/sh

echo "Stopping [[tft.deb.name]]... "
systemctl stop [[tft.deb.name]]
systemctl disable [[tft.deb.name]]
