# modules/Installers/LinuxServiceDEB/src/deb/control/prerm
#!/bin/sh

echo "Stopping bluemarine2-service... "
systemctl stop bluemarine2-service
systemctl disable bluemarine2-service
